# docker/supervisor/conf.d/queues.conf
; Email service queue workers — one program per priority queue.
; Scale the entire queue container: docker compose up -d --scale queue=3

[program:queue-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s --queue=emails-high --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
numprocs=%(ENV_QUEUE_WORKERS_HIGH)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s --queue=emails-default --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
numprocs=%(ENV_QUEUE_WORKERS_DEFAULT)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s --queue=emails-low --sleep=5 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
numprocs=%(ENV_QUEUE_WORKERS_LOW)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-bulk]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s --queue=emails-bulk --sleep=5 --tries=3 --timeout=180 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
numprocs=%(ENV_QUEUE_WORKERS_BULK)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:queue-retry]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s --queue=emails-retry --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
numprocs=%(ENV_QUEUE_WORKERS_RETRY)s
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
