# docker/supervisor/conf.d/queues.conf
; Email service queue workers — one program per priority queue.
; Scale the entire queue container: docker compose up -d --scale queue=3

[program:queue-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=emails-high --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=%(ENV_QUEUE_WORKERS_HIGH)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-high.log

[program:queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=emails-default --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=%(ENV_QUEUE_WORKERS_DEFAULT)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-default.log

[program:queue-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=emails-low --sleep=5 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=%(ENV_QUEUE_WORKERS_LOW)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-low.log

[program:queue-bulk]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=emails-bulk --sleep=5 --tries=3 --timeout=180 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=%(ENV_QUEUE_WORKERS_BULK)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-bulk.log

[program:queue-retry]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=emails-retry --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=%(ENV_QUEUE_WORKERS_RETRY)s
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-retry.log
